# sound_monitor.service
[Unit]
Description=Sound Level Monitor (Home Assistant MQTT)
After=network-online.target
Wants=network-online.target
# Stop retrying if the service crash-loops (more than 5 starts in 120s).
# These are [Unit] directives — systemd ignores them under [Service].
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
# On non-Raspberry-Pi hosts the `pi` user and /home/pi may not exist; edit
# User, WorkingDirectory, and the venv path in ExecStart to match your host
# (a mismatch here fails at startup with status=217/USER).
User=pi
WorkingDirectory=/home/pi
# Runs from a virtualenv — Pi OS forbids system-wide pip installs (PEP 668).
# Create it once with:
#   python3 -m venv /home/pi/venv
#   /home/pi/venv/bin/pip install -r /home/pi/requirements.txt
# Disables hardware AGC and pins capture gain before every start — see
# fix_mic_gain.py. A no-op if audio.capture_volume_percent isn't set in
# config.yaml. Runs with the same User/venv as ExecStart below.
ExecStartPre=/home/pi/venv/bin/python /home/pi/fix_mic_gain.py
ExecStart=/home/pi/venv/bin/python /home/pi/sound_monitor.py
Restart=always
RestartSec=10

# Keep stdout/stderr in the journal:
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sound_monitor

[Install]
WantedBy=multi-user.target
